# home-control-backend.service
[Unit]
Description=Home Control Backend (Django + Gunicorn)
Documentation=https://github.com/your-repo/home_control_adv
After=network.target
Wants=network.target

[Service]
Type=exec
User=manu
Group=manu
WorkingDirectory=/home/manu/code/home_control_adv
ExecStart=/home/manu/code/home_control_adv/start_backend.sh
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=15
TimeoutStartSec=30
PrivateTmp=true
Restart=always
RestartSec=10

# Variables de entorno
Environment=DJANGO_SETTINGS_MODULE=home_control.settings
Environment=PYTHONPATH=/home/manu/code/home_control_adv/backend
Environment=PYTHONDONTWRITEBYTECODE=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=home-control-backend

# Security (opcional, para mayor seguridad)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/manu/code/home_control_adv

[Install]
WantedBy=multi-user.target