# home-control-mqtt-bridge.service
[Unit]
Description=Home Control MQTT Bridge
After=network.target mosquitto.service home-control-django.service
Wants=mosquitto.service home-control-django.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/home_control
Environment=PATH=/home/pi/home_control/venv/bin
EnvironmentFile=/home/pi/home_control/.env
ExecStart=/home/pi/home_control/venv/bin/python mqtt_bridge.py
Restart=always
RestartSec=15

# Configuración de logs
StandardOutput=append:/var/log/home_control/mqtt_bridge.log
StandardError=append:/var/log/home_control/mqtt_bridge_error.log

# Configuración de seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/home_control
ReadWritePaths=/var/log/home_control

# Configuración de reinicio inteligente
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target